# app/src/main/java/Adobe/AdobeGCClient/com.adobe.agsservice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>Adobe_Genuine_Software_Integrity_Service</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/Adobe/AdobeGCClient/AGSService</string>
	</array>
	<key>Nice</key>
	<integer>1</integer>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/AlTest1.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/AlTest1.out</string>
</dict>
</plist>
